# systemd/.config/systemd/user/dotfiles-sync.service
[Unit]
Description=Pull dotfiles from remote

[Service]
Type=oneshot
ExecStart=%h/zphr-linux-config/scripts/dotfiles-sync.sh
Environment=DISPLAY=:0
