# systemd/.config/systemd/user/zphr-second-brain-backup.service
[Unit]
Description=Backup zphr-second-brain vault

[Service]
Type=oneshot
ExecStart=%h/zphr-linux-config/scripts/zphr-second-brain-backup.sh
