# systemd/.config/systemd/user/zphr-second-brain-backup.timer
[Unit]
Description=Hourly backup of zphr-second-brain

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
